# Da-HP/DA-HP-ISO/airootfs/etc/systemd/system/cpupower.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cpupower frequency-set -g performance
