# archive/scripts/aura-preload.service
[Unit]
Description=Aura Model Pre-Loader
After=ollama.service
Requires=ollama.service

[Service]
Type=oneshot
User=daripper
ExecStart=/usr/local/bin/aura-preload.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
